# app/pkg/s6/nats/run
#!/command/with-contenv sh
# Start NATS server

echo "Starting NATS server..."
exec nats-server -c /opt/streamops/app/pkg/nats.conf