# app/pkg/services.d/nats/run
#!/command/with-contenv sh

echo "Starting NATS server..."
exec nats-server -c /opt/streamops/app/pkg/nats.conf